# backup/services/readyroom/readyroom-smoke-test.timer
[Unit]
Description=Run ReadyRoom smoke test every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=readyroom-smoke-test.service

[Install]
WantedBy=timers.target
