# backup/services/m5/readyroom-smoke-test.service
[Unit]
Description=ReadyRoom smoke test

[Service]
Type=oneshot
ExecStart=/usr/local/bin/readyroom-smoke-test
