# backup/services/m5/spot-worker6.service
[Unit]
Description=M5 Worker 6 (Quadro M4000)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/spot-worker6
Environment=EXEC_TOKEN=changeme6
Environment=WORKER_NAME=m5_exec
Environment=CUDA_VISIBLE_DEVICES=0
Environment=JOBS_ROOT=/mnt/collective/jobs
Environment=JOBS_OUT_ROOT=/mnt/collective/jobs_out
ExecStart=/opt/spot-worker6/.venv/bin/uvicorn app:APP --host 0.0.0.0 --port 8789
Restart=always
RestartSec=2
User=ogre
Group=ogre

[Install]
WantedBy=multi-user.target
