# watch/recovery/spot-primary-lease-renew.timer
[Unit]
Description=Renew Spot Primary Witness Lease Every Fifteen Seconds

[Timer]
OnBootSec=10s
OnUnitActiveSec=15s
AccuracySec=1s
Persistent=true
Unit=spot-primary-lease-renew.service

[Install]
WantedBy=timers.target
