# backup/services/spot/spotd.service.d/20-worker-bias.conf
[Service]
Environment=WORKER_SCORE_BIAS_M5_EXEC=3000
Environment=WORKER_SCORE_BIAS_M5_WATCH=900
Environment=WORKER_SCORE_BIAS_SPOT_EXEC=0
Environment=SPOT_RESERVE_MIN_VRAM=10000
