# watch/remediation/spot-remediation-fixture.service
[Unit]
Description=Spot Remediation Fixture Service
Documentation=internal:starfleet-phase-2.38
ConditionPathExists=/usr/local/lib/spot/spot-remediation-fixture.sh

[Service]
Type=oneshot
ExecStart=/usr/local/lib/spot/spot-remediation-fixture.sh
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/spot/remediation-fixture
CapabilityBoundingSet=
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
