# backup/services/daystrom/spot-worker8.service
[Unit]
Description=M5 Worker 8 (GTX 1060)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/spot-worker8
Environment=EXEC_TOKEN=changeme8
Environment=WORKER_NAME=m5_watch
Environment=CUDA_VISIBLE_DEVICES=1
Environment=JOBS_ROOT=/mnt/collective/jobs
Environment=JOBS_OUT_ROOT=/mnt/collective/jobs_out
ExecStart=/opt/spot-worker8/.venv/bin/uvicorn app:APP --host 0.0.0.0 --port 8791
Restart=always
RestartSec=2
User=ogre
Group=ogre

[Install]
WantedBy=multi-user.target
