# backup/services/daystrom/spot-worker2.service
[Unit]
Description=Spot Worker 2 (Daystrom GPU)
After=network.target

[Service]
User=ogre
WorkingDirectory=/opt/spot-worker2
ExecStart=/opt/spot-worker2/.venv/bin/uvicorn app:APP --host 0.0.0.0 --port 8782
Restart=always

[Install]
WantedBy=multi-user.target
