# backup/services/spot/spot-dispatchd.service
[Unit]
Description=Spot Dispatcher Daemon (job queue -> remote agents)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ogre
Group=ogre
WorkingDirectory=/opt/spot/dispatch
ExecStart=/usr/bin/env python3 /opt/spot/dispatch/dispatchd.py
Restart=always
RestartSec=2
Environment=SPOT_DISPATCH_POLL_SECS=1.0
Environment=SPOT_DISPATCH_TIMEOUT=30

[Install]
WantedBy=multi-user.target
