# backup/services/spot/spot-gateway-smoke-test.service
[Unit]
Description=Spot Gateway smoke test

[Service]
Type=oneshot
ExecStart=/usr/local/bin/spot-gateway-smoke-test
