# backup/services/spot/spotd.service
[Unit]
Description=Spot Brain API (spotd)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/spotd.env
WorkingDirectory=/opt/spotd
ExecStart=/opt/spotd/.venv/bin/uvicorn app:APP --host 0.0.0.0 --port 8797 --workers 1 --loop asyncio --http h11 --log-level info
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target