# backup/services/spot/spot-helper.service
[Unit]
Description=Spot Helper (worker)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment="EXEC_TOKEN=changeme8"
WorkingDirectory=/opt/spot-helper
ExecStart=/opt/spot-helper/.venv/bin/uvicorn app:APP --host 0.0.0.0 --port 8791
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
