# backup/services/spot/spot-worker0.service.d/20-allowed-jobs.conf
[Service]
Environment=ALLOWED_JOBS=example_echo,example_sleep,llm_infer,llm_infer_stream
