# backup/services/m5/spot-worker8.service.d/gpu.conf
[Service]
Environment=CUDA_DEVICE_ORDER=PCI_BUS_ID
Environment=CUDA_VISIBLE_DEVICES=1
