# backup/services/spot/ollama.service.d/override.conf
[Service]
Environment="OLLAMA_KEEP_ALIVE=2m"
